# systemd/btrfs-backup.service
[Unit]
Description=Btrfs Automatic Backup Service
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/btrfs-backup.sh
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
